# tools/researcher-vm/deploy/systemd/openrx-cron@.service
[Unit]
Description=Run OpenRx background cron job %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=openrx
Group=openrx
UMask=0077
WorkingDirectory=/opt/openrx/current/tools/researcher-vm
EnvironmentFile=/etc/openrx/research.env
ExecStart=/opt/openrx/current/tools/researcher-vm/scripts/run-openrx-cron.sh %i

[Install]
WantedBy=multi-user.target
