# tools/researcher-vm/deploy/systemd/tao-trader-research.service
[Unit]
Description=Run one TAO autoresearch cycle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=taoresearch
Group=taoresearch
WorkingDirectory=/opt/tao-trader/current
EnvironmentFile=/etc/tao-trader/research.env
ExecStart=/opt/tao-trader/current/scripts/run-research-cycle.sh
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/tao-trader
